# ansible/roles/slurm_install/templates/slurm-node-maintenance.timer.j2
[Unit]
Description=Execute %N.service daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=4h

[Install]
WantedBy=timers.target
